# deploy/systemd/evaluate-venues.timer
[Unit]
Description=Run evaluate-venues daily at 06:00 UTC

[Timer]
OnCalendar=*-*-* 06:00:00 UTC
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
